# config/vault/supervisor.conf
[program:vault_server]
command=vault server -config=/etc/vault/vault.hcl 
autostart=true
autorestart=true
stderr_logfile=/var/log/vault.error.log
stdout_logfile=/var/log/vault.log
